# minecraftd-restart.service.in
[Unit]
Description=Daily Gaia Forge Server Restart
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/@INAME@ restart
User=@GAME_USER@

[Install]
WantedBy=multi-user.target
