# minecraftd-restart.timer.in
[Unit]
Description=Daily Gaia Forge Server Restart

[Timer]
OnCalendar=*-*-* 06:00:00
AccuracySec=5min
Persistent=true

[Install]
WantedBy=multi-user.target
